# apache-server/src/supervisor/supervisord.apache.conf
[program:apache]
command=httpd -f /etc/apache2/httpd.conf -DFOREGROUND
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
